# supervisord/supervisord.conf
[supervisord]
logfile=/var/log/supervisor.log ;
logfile_maxbytes=50MB        ;
logfile_backups=10           ;
loglevel=info                ;
pidfile=/tmp/supervisord.pid ;
nodaemon=false               ;
minfds=1024                  ;
minprocs=200                 ;


[program:fireapps]
command=/deploy/envs/sales/bin/gunicorn conf.wsgi:application --bind 127.0.0.1:7007 --pid /tmp/app.pid ;
directory=/deploy/apps/fire_apps/ ;

[program:worker]
command=/deploy/envs/sales/bin/celery --app=conf.celery:app worker --loglevel=INFO ;
directory=/deploy/apps/fire_apps/ ;

[program:beat]
command=/deploy/envs/sales/bin/celery  --app=conf.celery:app beat --loglevel=INFO -S django ;
directory=/deploy/apps/fire_apps/ ;

[inet_http_server]
port=127.0.0.1:9001   ;

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface ;

[supervisorctl]
serverurl=http://127.0.0.1:9001 ;
